# mongo.service.conf
[program:mongo]
command=/usr/bin/mongod --rest
